# examples/launchd/com.owen.ai-flow.tick.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.owen.ai-flow.tick</string>
    <key>ProgramArguments</key>
    <array>
      <string>REPLACE_WITH_NODE_PATH</string>
      <string>REPLACE_WITH_AI_FLOW_SCRIPT_PATH</string>
      <string>scan</string>
    </array>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/ai-flow.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/ai-flow.err.log</string>
  </dict>
</plist>
